# ai.montageai.cookies.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.montageai.cookies</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/omaligamini4/Desktop/montageai/refresh_cookies.sh</string>
    </array>

    <!-- Corre cada día a las 6:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- También corre al cargar el plist (primera vez / después de reboot) -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/omaligamini4/Library/Logs/montageai_cookies.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/omaligamini4/Library/Logs/montageai_cookies.log</string>

    <!-- Si falla, reintentar en 60 segundos -->
    <key>StartInterval</key>
    <integer>0</integer>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
